# templates/apprc_service.epp
# chkconfig: 0136 70 39 
### BEGIN INIT INFO
# Provides: APPSVC 
# Required-Start:
# Required-Stop:
# Default-Start: 3
# Default-Stop: 0 1 6 
### END INIT INFO

case "$1" in
    
start)
     echo starting script
     ;;
stop)
    echo stoping script
     ;;
restart)
    INFO "Restart is not supported by $0"
     ;;

 *)
    INFO "Usage: $0 {start|stop}"
     ;;

esac

exit 0

######### THIS IS THE END OF THE STARTUP SCRIPT TEMPLATE ########